# systemd/overtime-kiosk.service
[Unit]
Description=Overtime Kiosk
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/tracymaint/overtime_pi_kiosk_full
Environment="PYTHONUNBUFFERED=1"
ExecStart=/home/tracymaint/overtime_pi_kiosk_full/.venv/bin/python app.py
Restart=always
User=tracymaint
Group=tracymaint

# Give the service access to display; adjust if using Wayland
Environment=DISPLAY=:0

[Install]
WantedBy=multi-user.target
