# systemd/overtime-kiosk.service
[Unit]
Description=Overtime Kiosk v0.9.0
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/overtime_pi_kiosk_full
Environment="FLASK_ENV=production"
ExecStart=%h/overtime_pi_kiosk_full/.venv/bin/python run.py
Restart=on-failure

[Install]
WantedBy=multi-user.target
